# scripts/run_blogs.cron
1 6 * * * /var/www/django-web/blog/scripts/run_blogs.sh >> crawl_blog.log
